# crontab
# Send out any 101 questions emails for this hour
5 * * * * root cd /app && . ./env.conf && python3 cron.py
